# fpm/recipes/prometheus-node-exporter-0.17.0/upstart-node-exporter.conf
# node-exporter: Prometheus Node Exporter - to make metrics available
#
# Homepage: https://prometheus.io
#

description "node-exporter: Prometheus Node Exporter - to make metrics available"

start on (local-filesystems and net-device-up IFACE!=lo)

respawn

script
/usr/local/bin/node_exporter --web.listen-address=:9080
end script
