# pkg/statsd/debian/statsd.upstart
# statsd - Network daemon for aggregating statistics
#
# This is a network service that receives metric data via UDP from other
# applications. It aggregates this data and flushes it to a storage backend
# (typically Graphite) at regular intervals.
#
description "StatsD: a simple daemon for easy stats aggregation"
author "Etsy"

start on (local-filesystems and net-device-up IFACE!=lo)

setuid statsd
setgid statsd

respawn
respawn limit 10 5

exec /usr/bin/statsd /etc/statsd.conf
